# pkg/v1/zcli/usr/lib/zyphor/pkg/install
#!/bin/bash

zyphor_pkg_install() {

  APP="$3"

  # ==================================================
  # VALIDATE INPUT
  # ==================================================
  if [ -z "$APP" ]; then
    echo ""
    echo "Zyphor: ✖ Missing app name"
    echo "Usage: zyphor pkg install <app>"
    echo ""
    return 1
  fi

  # ==================================================
  # ROOT CHECK
  # ==================================================
  if [ "$EUID" -ne 0 ]; then
    echo ""
    echo "Zyphor: ⚠ Requires root privileges"
    echo "Try: sudo zyphor pkg install $APP"
    echo ""
    return 1
  fi

  # ==================================================
  # FETCH REGISTRY
  # ==================================================
  echo ""
  echo "Zyphor: Fetching registry..."

  if ! curl -fsSL "$REGISTRY_URL" -o "$CACHE_FILE"; then
    echo "Zyphor: ⚠ Failed to fetch remote registry"

    if [ ! -f "$CACHE_FILE" ]; then
      echo "Zyphor: ✖ No local cache available"
      return 1
    fi

    echo "Zyphor: Using cached registry..."
  fi

  # ==================================================
  # RESOLVE APP
  # ==================================================
  TYPE=$(jq -r ".apps.\"$APP\".type" "$CACHE_FILE")
  PACKAGE=$(jq -r ".apps.\"$APP\".package // empty" "$CACHE_FILE")
  URL=$(jq -r ".apps.\"$APP\".url // empty" "$CACHE_FILE")
  MODULE=$(jq -r ".apps.\"$APP\".module // empty" "$CACHE_FILE")

  FLAGS=$(jq -r ".apps.\"$APP\".flags // empty" "$CACHE_FILE")
  CHANNEL=$(jq -r ".apps.\"$APP\".channel // empty" "$CACHE_FILE")

  # Optional checksum
  SHA256=$(jq -r ".apps.\"$APP\".sha256 // empty" "$CACHE_FILE")

  # ==================================================
  # VALIDATE APP
  # ==================================================
  if [ "$TYPE" = "null" ] || [ -z "$TYPE" ]; then
    echo ""
    echo "Zyphor: ✖ App not found → $APP"
    echo ""
    return 1
  fi

  echo ""
  echo "Zyphor: Installing '$APP' (type: $TYPE)"
  echo ""

  # ==================================================
  # INSTALL HANDLER
  # ==================================================
  case "$TYPE" in

    # --------------------------
    # APT INSTALL
    # --------------------------
    apt)
      if [ -z "$PACKAGE" ]; then
        echo "Zyphor: ✖ Missing package name"
        return 1
      fi

      echo "Zyphor: APT install → $PACKAGE"
      apt update -qq
      apt install -y $FLAGS "$PACKAGE"
      ;;

    # --------------------------
    # DEB INSTALL
    # --------------------------
    deb)
      if [ -z "$URL" ]; then
        echo "Zyphor: ✖ Missing .deb URL"
        return 1
      fi

      TMP_DEB="/tmp/zyphor_${APP}.deb"

      echo "Zyphor: Downloading .deb → $URL"
      if ! wget -q --show-progress "$URL" -O "$TMP_DEB"; then
        echo "Zyphor: ✖ Failed to download .deb"
        return 1
      fi

      # ----- CHECKSUM VALIDATION (OPTIONAL) -----
      if [ -n "$SHA256" ] && [ "$SHA256" != "null" ]; then
        echo "Zyphor: Verifying checksum..."
        DOWNLOADED_HASH=$(sha256sum "$TMP_DEB" | awk '{print $1}')

        if [ "$DOWNLOADED_HASH" != "$SHA256" ]; then
          echo "Zyphor: ✖ Checksum mismatch!"
          rm -f "$TMP_DEB"
          return 1
        fi
      fi

      echo "Zyphor: Installing .deb..."
      dpkg -i $FLAGS "$TMP_DEB"

      echo "Zyphor: Fixing dependencies..."
      apt -f install -y

      echo "Zyphor: Cleaning up..."
      rm -f "$TMP_DEB"
      ;;

      # --------------------------
      # SNAP INSTALL (FIXED + BOOTSTRAP SAFE)
      # --------------------------
      snap)

        if [ -z "$PACKAGE" ]; then
          echo "Zyphor: ✖ Missing snap package"
          return 1
        fi

        echo "Zyphor: Snap install → $PACKAGE"

        # ==================================================
        # INSTALL SNAPD IF MISSING
        # ==================================================
        if ! command -v snap >/dev/null 2>&1; then
          echo "Zyphor: snapd not found → installing..."

          apt update -qq
          apt install -y snapd
        fi

        # ==================================================
        # ENSURE SNAP + SECURITY SERVICES
        # ==================================================
        echo "Zyphor: Initializing snapd service..."

        systemctl enable --now snapd >/dev/null 2>&1
        systemctl enable --now snapd.socket >/dev/null 2>&1

        # Enable AppArmor for snap confinement (NEW)
        systemctl enable --now snapd.apparmor >/dev/null 2>&1

        # (Optional but strong baseline)
        systemctl enable --now apparmor >/dev/null 2>&1

        # Wait for daemon readiness (VERY IMPORTANT)
        echo "Zyphor: Waiting for snapd to be ready..."
        sleep 3

        # ==================================================
        # DEFAULT CHANNEL
        # ==================================================
        if [ -z "$CHANNEL" ] || [ "$CHANNEL" = "null" ]; then
          CHANNEL="stable"
        fi

        echo "Zyphor: Channel → $CHANNEL"

        # ==================================================
        # INSTALL SNAP PACKAGE (RETRY SAFE)
        # ==================================================
        snap install "$PACKAGE" $FLAGS --channel="$CHANNEL"
        SNAP_EXIT=$?

        # ==================================================
        # RETRY ON FIRST-BOOT FAILURE (common issue)
        # ==================================================
        if [ "$SNAP_EXIT" -ne 0 ]; then
          echo "Zyphor: snap install failed → retrying after refresh..."

          systemctl restart snapd >/dev/null 2>&1
          sleep 2

          snap install "$PACKAGE" $FLAGS --channel="$CHANNEL"
          SNAP_EXIT=$?
        fi

        return $SNAP_EXIT
        ;;
    # --------------------------
    # INTERNAL MODULE
    # --------------------------
    internal)
      if [ -z "$MODULE" ]; then
        echo "Zyphor: ✖ Missing module"
        return 1
      fi

      echo "Zyphor: Running module → $MODULE"
      bash "/usr/lib/zyphor/modules/$MODULE/install.sh"
      ;;

    # --------------------------
    # UNKNOWN TYPE
    # --------------------------
    *)
      echo "Zyphor: ✖ Unknown type → $TYPE"
      return 1
      ;;
  esac

  # ==================================================
  # RESULT
  # ==================================================
  if [ $? -eq 0 ]; then
    echo ""
    echo "Zyphor: ✔ Installed '$APP' successfully"
    echo ""
  else
    echo ""
    echo "Zyphor: ✖ Failed installing '$APP'"
    echo ""
    return 1
  fi
}