# os-config/variant-xfce/includes.chroot/etc/skel/.config/autostart/install.desktop
[Desktop Entry]
Type=Application
Name=Install Zyphor
Exec=calamares
Icon=system-installer
X-GNOME-Autostart-enabled=true